# root/etc/uci-defaults/30_luci-theme-m3e
#!/bin/sh

# Keep M3E in a separate namespace while reusing existing material3 static directories.
# Use option names without '-' for better UCI compatibility.
uci set luci.themes.M3E=/luci-static/material3
uci set luci.themes.M3EBlue=/luci-static/material3-blue
uci set luci.themes.M3EGreen=/luci-static/material3-green
uci set luci.themes.M3ERed=/luci-static/material3-red

# Clean up legacy keys from previous naming experiments.
uci -q delete luci.themes.M3E-Blue
uci -q delete luci.themes.M3E-Green
uci -q delete luci.themes.M3E-Red

if [ "$PKG_UPGRADE" != 1 ]; then
	uci set luci.main.mediaurlbase=/luci-static/material3
fi

uci commit luci

exit 0